# cron/send_email.cron
* * * * * python /www/ischoolbus/manage.py send_mail --settings=ischoolbus.production
